# systemd/bridge-worker@.service
[Unit]
Description=ctrl+shft bridge worker (instance %i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/dotfiles/secrets/.env.agent
EnvironmentFile=%h/dotfiles/secrets/.env.secrets
Environment=WORKER_ID=%i
WorkingDirectory=%h/dotfiles
ExecStart=%h/dotfiles/secrets/.venv/bin/python -m bridge.worker
Restart=on-failure
RestartSec=10
StandardOutput=append:%h/bridge/logs/bridge.log
StandardError=append:%h/bridge/logs/bridge.log

# Hardening
NoNewPrivileges=yes
# MVP: worker needs access to both ~/bridge and ~/dotfiles/working
# Phase 2: tighten to just ~/bridge once per-workspace working/ lands
ReadWritePaths=%h/bridge %h/dotfiles/working

[Install]
WantedBy=default.target
